# metric-collector/proxy-exporter.service
# Ubuntu services - proxy client - PUB-2147063-2310182036
[Unit]
Description=Start the Metrics Collector v5 proxy client

[Service]
Type=simple
User=root
Restart=always
RestartSec=10
CPUAccounting=yes
MemoryAccounting=yes

ExecStart=/opt/eyeflow/monitor/lib/pushprox-client --proxy-url=https://proxy.customer.eyeflow.ai/ --log.level=warn

[Install]
WantedBy=multi-user.target

