# metric-collector/proxy-exporter.service
# Ubuntu services - proxy client - PUB-46afabb-20230816125700
[Unit]
Description=Start the Metrics Collector v5 proxy client

[Service]
Type=simple
User=eyeflow
Restart=always
RestartSec=10
CPUAccounting=yes
MemoryAccounting=yes

ExecStart=/opt/eyeflow/monitor/lib/pushprox-client --proxy-url=http://34.28.1.197:9096/

[Install]
WantedBy=multi-user.target

