# metric-collector/proxy-exporter.service
# Ubuntu services - proxy client - PUB-535e5c0-20230817194124
[Unit]
Description=Start the Metrics Collector v5 proxy client

[Service]
Type=simple
User=eyeflow
Restart=always
RestartSec=10
CPUAccounting=yes
MemoryAccounting=yes

ExecStart=/opt/eyeflow/monitor/lib/pushprox-client --proxy-url=https://proxy.monitoring.eyeflow.ai/

[Install]
WantedBy=multi-user.target

